# packaging/martinspoint.service
[Unit]
Description=Martins Point Security
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/martinspoint
ExecStart=/opt/martinspoint/run_idview.sh
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1
Environment=QT_QPA_PLATFORM=xcb

[Install]
WantedBy=multi-user.target
